# rootfs/usr/lib/opkg/info/firewall.list
/lib/firewall/uci_firewall.sh
/etc/config/firewall
/etc/init.d/firewall
/etc/hotplug.d/iface/20-firewall
/etc/firewall.user
