# rootfs/usr/lib/opkg/info/luci-admin-core.list
/etc/init.d/luci_fixtime
/etc/init.d/luci_dhcp_migrate
/lib/uci/schema/default/fstab
/lib/uci/schema/default/httpd
/lib/uci/schema/default/system
/lib/uci/schema/default/network
/lib/uci/schema/default/luci_hosts
/lib/uci/schema/default/wireless
/lib/uci/schema/default/dhcp
/lib/uci/schema/default/dropbear
/sbin/luci-flash
/usr/lib/lua/luci/tools/webadmin.lua
/usr/lib/lua/luci/view/header.htm
/usr/lib/lua/luci/view/about.htm
/usr/lib/lua/luci/view/error500.htm
/usr/lib/lua/luci/view/footer.htm
/usr/lib/lua/luci/view/indexer.htm
/usr/lib/lua/luci/view/sysauth.htm
/usr/lib/lua/luci/view/error404.htm
