# rootfs/usr/lib/opkg/info/uhttpd.list
/etc/init.d/uhttpd
/etc/config/uhttpd
/usr/sbin/uhttpd
